# projects/backups/backup-config-files/stpi-proxy.virtual-labs.ac.in/vz/private/12237/etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

#Add default gw as stpi-router.virtual-labs.ac.in for Internet access
route add default gw 10.4.12.236

#Add routes for local networks via Local L3 Gateway
route add -net 10.0.0.0/8 gw 10.4.12.1
route add -net 172.16.0.0/12 gw 10.4.12.1
route add -net 192.168.0.0/16 gw 10.4.12.1

